# root/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/sbin/openrc shutdown

tty1::respawn:/usr/local/bin/rustfs-dashboard
ttyS0::respawn:/usr/local/bin/rustfs-dashboard
