# root/etc/init.d/rustfs
#!/sbin/openrc-run
description="RustFS S3-compatible object storage"

depend() {
    need net alfpine-disks
    after localmount
}

start() {
    ebegin "Starting RustFS"
    export RUSTFS_ROOT_USER="${RUSTFS_ROOT_USER:-rustfsadmin}"
    export RUSTFS_ROOT_PASSWORD="${RUSTFS_ROOT_PASSWORD:-rustfsadmin}"
    export RUSTFS_ADDRESS="${RUSTFS_ADDRESS:-:9000}"
    export RUSTFS_CONSOLE_ADDRESS="${RUSTFS_CONSOLE_ADDRESS:-:9001}"

    VOLUMES=$(cat /run/alfpine-volumes 2>/dev/null)
    if [ -z "$VOLUMES" ]; then
        VOLUMES="/data"
    fi

    start-stop-daemon --start --background --make-pidfile \
        --pidfile /run/rustfs.pid \
        --stdout /var/log/rustfs.log --stderr /var/log/rustfs.log \
        --exec /usr/local/bin/rustfs -- server $VOLUMES
    eend $?
}

stop() {
    ebegin "Stopping RustFS"
    start-stop-daemon --stop --pidfile /run/rustfs.pid
    eend $?
}
